# torrserver/rootfs/etc/s6-overlay/s6-rc.d/torserver/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: torserver
# Runs the motionEye
# ==============================================================================
declare -a options

bashio::log.info "Starting torserver..."


# Enable debug mode
if bashio::debug; then
    options+=(-d)
fi

# Run the torserver
exec meyectl startserver "1"
